# infra/systemd/leafora-api.service
[Unit]
Description=Leafora API
After=network.target postgresql.service

[Service]
User=leafora
Group=leafora
WorkingDirectory=/var/www/leafora/apps/api
EnvironmentFile=/etc/leafora/api.env
ExecStart=/var/www/leafora/apps/api/.venv/bin/uvicorn app.main:create_app --factory --host 127.0.0.1 --port 18473
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
